# com.nicolassodoski.processar-mensagens.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nicolassodoski.processar-mensagens</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/nicolassodoski/Claude/Projects/PW/processar_e_push.sh</string>
    </array>

    <!-- Rodar a cada 30 minutos (1800 segundos) -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/nicolassodoski/Claude/Projects/PW/processar_auto.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/nicolassodoski/Claude/Projects/PW/processar_auto.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/nicolassodoski</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
