# com.imoveis.verificador-erros.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.imoveis.verificador-erros</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/nicolassodoski/Claude/Projects/PW/verificador_erros.py</string>
    </array>

    <!-- Rodar a cada 1 hora (3600 segundos) -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/nicolassodoski/Claude/Projects/PW/verificador_erros_log.txt</string>

    <key>StandardErrorPath</key>
    <string>/Users/nicolassodoski/Claude/Projects/PW/verificador_erros_log.txt</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/nicolassodoski</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
